# teoria/02-iac/02-network-debugging/david-canedo/files/network-watchdog.timer
[Unit]
Description=Network Watchdog Timer

[Timer]
OnBootSec=5
OnUnitActiveSec=5s
Unit=network-watchdog.service

[Install]
WantedBy=timers.target

